# systemd/oracle-bot.service
[Unit]
Description=Delphi-42 Meshtastic oracle bot
After=local-fs.target llm-openai-api.service
Wants=local-fs.target llm-openai-api.service

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/delphi-42
Environment=DELPHI_CONFIG=/opt/delphi-42/config/oracle.example.yaml
Environment=DELPHI_LOG_LEVEL=INFO
ExecStart=/usr/bin/env python -m bot.oracle_bot
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
