# systemd/oracle-core.service
[Unit]
Description=Delphi-42 offline index refresh
After=local-fs.target
Requires=local-fs.target

[Service]
Type=oneshot
User=pi
WorkingDirectory=/opt/delphi-42
ExecStart=/usr/bin/env python -m ingest.build_index --input-dir /opt/delphi-42/data/library/plaintext --db /opt/delphi-42/data/index/oracle.db

[Install]
WantedBy=multi-user.target
